# systemctl-timers/kdata.timer
[Unit]
Description=Daily kdata generation

[Timer]
OnCalendar=*-*-* 17:30:00

[Install]
WantedBy=multi-user.target
RequiredBy=multi-user.target

